# roles/nextcloud/templates/notify_push.supervisord.ini.j2
[program:notify_push]
command=notify_push {{ uberspace_nextcloud__document_root }}/config/config.php
autostart=yes
autorestart=yes
startsecs=30
