# roles/mosquitto/templates/supervisorctl.ini.j2
[program:mosquitto-{{ uberspace_mosquitto__instance_name }}]
command=mosquitto -c %(ENV_HOME)s/etc/mosquitto/mosquitto-{{ uberspace_mosquitto__instance_name }}.conf
autostart=yes
autorestart=yes
